# workspace/scripts/files/system/systemd/lighttpd.service
[Unit]
Description=Lighttpd Daemon
After=network-online.target am_i_leader.service
Requires=am_i_leader.service

[Service]
Type=simple
PIDFile=/run/lighttpd.pid
ExecStartPre=/usr/sbin/lighttpd -tt -f /etc/lighttpd/lighttpd.conf
ExecStart=/usr/sbin/lighttpd -D -f /etc/lighttpd/lighttpd.conf
ExecReload=/bin/kill -USR1 $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target